# conf-default/cron.d/nmis9
# NMIS9 Cron Configuration

# if you DON'T want any NMIS cron mails to go to root,
# uncomment and adjust the next line
# MAILTO=WhoeverYouAre@yourdomain.tld

# some tools like fping reside outside the minimal path
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

# note that almost all operations are now performed by the nmisd

# m h dom month dow user command

# create reports, Daily Monthly Weekly
9 0 * * *  root /usr/local/nmis9/bin/nmis-cli act=run-reports period=day type=all
9 1 * * 0  root /usr/local/nmis9/bin/nmis-cli act=run-reports period=week type=all
9 2 1 * *  root /usr/local/nmis9/bin/nmis-cli act=run-reports period=month type=all

# performance data
*/5 * * * * root /usr/local/nmis9/bin/nmis-cli act=collect-top-data > /dev/null
0 * * * *  root /usr/local/nmis9/bin/nmis-cli act=collect-performance-data > /dev/null